# quadlets/media-jellyfin.container
[Unit]
Description=Jellyfin Media Server
After=network-online.target

[Container]
ContainerName=jellyfin
Image=docker.io/jellyfin/jellyfin:latest
UserNS=keep-id

# Keep memory constrained; Jellyfin can get hungry during library scans
Memory=1G

# Render devices for Intel QuickSync hardware transcoding
AddDevice=/dev/dri:/dev/dri
Annotation=run.oci.keep_groups=true

# SELinux: :Z for private config, :z for shared media directory
Volume=%h/srv/@media/config/jellyfin:/config:Z
Volume=%h/srv/@media/data:/data:z
Network=media-internal.network

# Only expose over Tailscale by binding to the Tailscale IP, or leave 
# as 0.0.0.0 if firewalld handles your Tailscale zone restriction.
PublishPort=8096:8096

[Install]
WantedBy=default.target

[Service]
CPUShares=1024
Restart=on-failure
RestartSec=5
